# static/etc/network/if-pre-up.d/iptables.sh
#!/bin/sh

cat /etc/sysconfig/iptables

/sbin/iptables-restore < /etc/sysconfig/iptables

