# static/var/spool/cron/crontabs/root
00   00 * * * /usr/sbin/rdate -s time.bora.net && /sbin/hwclock -w
00-59 * * * * check_taskmon &
